# apps/nexus-config/config/supervisor-nexus-chown-data.conf
[program:nexus-chown-data]
command=/bin/bash -xc "chown nexus:nexus /nexus-data"
user=root
autostart=false
startsecs=0
autorestart=0
startretries=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
